# mkosi.extra/usr/lib/systemd/system/getty@.service.d/10-enable-numlock.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

[Service]
ExecStartPre=/bin/sh -c 'setleds -D +num < /dev/%I'
